# hosting-utils/crontab
*/10 * * * * /var/www/html/LiveActionMap/update.sh
